# deploy/mac/com.aria.mlx-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aria.mlx-server</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3.14</string>
        <string>-m</string>
        <string>mlx_lm.server</string>
        <string>--model</string>
        <!-- ACTIVE: Qwen3.5-4B (2.9GB RAM, lower swap pressure) -->
        <string>mlx-community/Qwen3.5-4B-MLX-4bit</string>
        <!-- OLD: Qwen3.5-9B (5.1GB RAM) -->
        <!-- <string>mlx-community/Qwen3.5-9B-MLX-4bit</string> -->
        <!-- OLD: Qwen3-4B (deprecated) -->
        <!-- <string>mlx-community/Qwen3-4B-Instruct-2507-4bit</string> -->
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8080</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/mlx_server.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/mlx_server.err</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/najia</string>
    </dict>
    
    <key>WorkingDirectory</key>
    <string>/Users/najia</string>
</dict>
</plist>
